# configs/systemd/user/sync-screenshots-google-photos-fedora.path
[Unit]
Description = Sync Fedora's Screenshot Directory

[Path]
PathModified = %h/Pictures/Screenshots
Unit = sync-screenshots-google-photos.service

[Install]
WantedBy = default.target
