# configs/systemd/user/sync-guru-issues.timer
[Unit]
Description = Sync Guru GitLab issues to Trello

[Timer]
OnCalendar = hourly
Persistent = true

[Install]
WantedBy = timers.target
