# configs/systemd/user/sync-guru-issues.service
[Unit]
Description = Sync Guru GitLab issues to Trello

[Service]
Type = oneshot
ExecStart = python3 %h/.config/ewconfig/scripts/guru-sync-issues